# crontab
# =============================================================================
# DATABASE BACKUP & SYNC AUTOMATION - GÜNLÜK YEDEK
# Server: <YOUR_DB_HOSTNAME>
# Last Updated: 2026-07-21
# =============================================================================
#
# NEDEN GÜNLÜK (15 dakikada bir DEĞİL):
# Yedekler DB container'ının İÇİNDE `docker exec` (mariadb-dump / mongodump) ile
# alınır; dump'ın bellek tüketimi DB container'ının KENDİ cgroup'una sayılır.
# 15 dakikada bir tam yedek, MariaDB'yi tekrar tekrar CONSTRAINT_MEMCG OOM'a
# soktu (container 172 kez restart etti). Günde bir tam yedek + offsite sync
# bu baskıyı ortadan kaldırır.
#
# TIMELINE (Her Gün):
# 02:00 → Full backup (tüm DB'ler)
# 02:30 → Google Drive sync (yedek bittikten sonra)
# 03:00 → Retention temizliği (RETENTION_DAYS=7, backup.sh clean)
# 08:00 → İstatistik raporu
# 03:30 (Pazar) → 7 günden eski logları temizle
#
# =============================================================================


# -----------------------------------------------------------------------------
# BACKUP - Her gün 02:00 (Full backup)
# -----------------------------------------------------------------------------
0 2 * * * /opt/databases/backup.sh all >> /opt/databases/logs/cron_backup.log 2>&1

# -----------------------------------------------------------------------------
# GOOGLE DRIVE SYNC - Yedek bittikten sonra (02:30)
# -----------------------------------------------------------------------------
30 2 * * * /opt/databases/sync_remote.sh >> /opt/databases/logs/cron_sync.log 2>&1

# -----------------------------------------------------------------------------
# MAINTENANCE - Temizlik & Raporlama
# -----------------------------------------------------------------------------
# Eski yedekleri temizle - Her gün 03:00. backup.sh clean, RETENTION_DAYS'i
# (varsayılan 7 gün) kullanır → tek yetki kaynağı; kod ve doküman artık uyumlu.
# (Eski satır `-mmin +360` = 6 saat siliyordu, belgelenen 7 günle çelişiyordu.)
0 3 * * * /opt/databases/backup.sh clean >> /opt/databases/logs/cron_clean.log 2>&1

# Günlük istatistik raporu - Her gün 08:00
0 8 * * * /opt/databases/backup.sh stats >> /opt/databases/logs/cron_stats.log 2>&1

# Log dosyalarını temizle - Her Pazar 03:30 (7 günden eski)
30 3 * * 0 find /opt/databases/logs -name "*.log" -mtime +7 -delete 2>/dev/null

# -----------------------------------------------------------------------------
# OPTIONAL: Disk space alert (uncomment to enable)
# -----------------------------------------------------------------------------
# 0 */6 * * * df -h /opt/databases | awk 'NR==2 && int($5) > 80 {print "DISK WARNING: "$5" used"}' | mail -s "Disk Alert" admin@example.com

# =============================================================================
# NOTES:
# =============================================================================
# - Backup + sync flock ile korunur; eşzamanlı çalışma olmaz (bkz. backup.sh).
# - Günde 4 tam yedek dosyası alınır (MariaDB + PostgreSQL + MongoDB + Redis).
# - 7 günlük retention = ~28 tam yedek dosyası (+ tekil/manuel yedekler).
# - Tahmini disk kullanımı: DB boyutuna göre değişir.
#
# MANUAL COMMANDS:
# - Manual backup:  /opt/databases/backup.sh all
# - Manual sync:    /opt/databases/sync_remote.sh
# - Check status:   /opt/databases/backup.sh stats
# - List backups:   /opt/databases/backup.sh list
# - Test sync:      /opt/databases/sync_remote.sh test
# =============================================================================
